# udev/99-tym.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="tym -v apply"
